# microsoft-tts/rootfs/etc/s6-overlay/s6-rc.d/microsoft/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start Microsoft service
# ==============================================================================
flags=()
if bashio::config.true 'update_voices'; then
    flags+=('--update-voices')
fi

if bashio::config.true 'debug_logging'; then
    flags+=('--debug')
fi

if bashio::config.exists 'rate'; then
    flags+=("--rate $(bashio::config 'rate')")
fi

if bashio::config.exists 'pitch'; then
    flags+=("--pitch $(bashio::config 'pitch')")
fi

if bashio::config.exists 'volume'; then
    flags+=("--volume $(bashio::config 'volume')")
fi

if bashio::config.exists 'style'; then
    flags+=("--style $(bashio::config 'style')")
fi

if bashio::config.exists 'style_degree'; then
    flags+=("--style-degree $(bashio::config 'style_degree')")
fi

exec uv run -m wyoming_microsoft_tts \
    --uri "tcp://0.0.0.0:10200" \
    --subscription-key "$(bashio::config 'subscription_key')" \
    --service-region "$(bashio::config 'service_region')" \
    --download-dir /data \
    ${flags[@]}
